# deploy/walaspro.service
[Unit]
Description=WALIKELASPRO - Aplikasi Manajemen Wali Kelas
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/walas.my.id
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=walaspro
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
